# docker/docker-lnp/conf/supervisord.conf
[unix_http_server]
file = /tmp/supervisor.sock

[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
#设置为前台启动
nodaemon = true
minfds = 1024
minprocs = 200
umask = 022
#切换用户为root
user = root
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
startretries=3
priority=10

[program:php-fpm]
command=/usr/sbin/php-fpm7
autostart=true
autorestart=true
startretries=3
priority=10