# docker-WebSphereMQ/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=su mqm -c /start_queue_manager.sh
